# src/deb/control/postinst
#!/bin/sh -e

# Enable service
systemctl enable mastodon-exporter.service

# Stop service if started
systemctl stop mastodon-exporter.service