# src/deb/control/prerm
#!/bin/sh -e

# Disable service
systemctl disable mastodon-exporter.service